# deploy/solara-backend.service
[Unit]
Description=Solara backend (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
# Replace with the deploy user/group and the path this repo is cloned to.
User=solara
Group=solara
WorkingDirectory=/opt/solara
Environment="ALLOWED_ORIGINS=https://solarasim.space"
ExecStart=/opt/solara/.venv/bin/uvicorn server.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

# de440t.bsp is opened once at startup and held for the process lifetime --
# nothing else here needs elevated privileges or filesystem access beyond
# this directory and its own exports/ subfolder.
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/opt/solara/exports

[Install]
WantedBy=multi-user.target
